# crontab

# NOESIS Daily Scanner - Runs at MIDNIGHT UTC
# To add to crontab: crontab -e
# Then add this line:

0 0 * * * /usr/bin/python3 /root/.openclaw/workspace/noesis-prediction/cron_scanner.py >> /root/.openclaw/workspace/noesis-prediction/cron_logs/scanner.log 2>&1

# Or use python directly:
# 0 0 * * * cd /root/.openclaw/workspace/noesis-prediction && python3 cron_scanner.py >> cron_logs/scanner.log 2>&1
0 9 * * * cd /root/.openclaw/workspace/noesis-prediction && python3 morning_scan.py >> cron_logs/morning_scan.log 2>&1
0 7 * * * cd /root/.openclaw/workspace/noesis-prediction && python3 results_scanner.py >> cron_logs/results.log 2>&1
0 8 * * * cd /root/.openclaw/workspace/noesis-prediction && python3 vegas_odds.py >> cron_logs/vegas.log 2>&1
